# b/386/Linux/x/systemd/system/inferno.service
[Unit]
Description=inferno

[Service]
Type=forking
User=inferno
ExecStart=/b/386/Linux/x/init/inferno

[Install]
WantedBy=multi-user.target

